# .config/systemd/user/restic-backup.service
[Unit]
Description=restic backup to Proton Drive via proton-drive-sync
After=network-online.target proton-drive-sync.service
Wants=network-online.target

[Service]
Type=oneshot
Nice=19
IOSchedulingClass=idle
ExecStart=%h/backup.sh
